# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:pulseaudio]
command=pulseaudio --start --load="module-null-sink sink_name=virtual_sink" --load="module-virtual-source source_name=virtual_source" --exit-idle-time=-1
autorestart=true
priority=5
user=root
stdout_logfile=/var/log/supervisor/pulseaudio.log
stderr_logfile=/var/log/supervisor/pulseaudio.err

[program:xvfb]
command=Xvfb :1 -screen 0 %(ENV_SCREEN_WIDTH)sx%(ENV_SCREEN_HEIGHT)sx%(ENV_SCREEN_DEPTH)s -ac +extension GLX
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.err

[program:openbox]
command=openbox --display :1
autorestart=true
priority=20
environment=DISPLAY=":1"
stdout_logfile=/var/log/supervisor/openbox.log
stderr_logfile=/var/log/supervisor/openbox.err

[program:x11vnc]
command=x11vnc -display :1 -forever -shared -nopw -listen 0.0.0.0 -port %(ENV_VNC_PORT)s
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.err

[program:novnc]
command=websockify --web /usr/share/novnc %(ENV_NOVNC_PORT)s localhost:%(ENV_VNC_PORT)s
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.err

[program:stream]
command=/app/./stream.sh
autorestart=true
priority=80
environment=DISPLAY=":1",SCREEN_WIDTH="%(ENV_SCREEN_WIDTH)s",SCREEN_HEIGHT="%(ENV_SCREEN_HEIGHT)s",STREAM_FPS="%(ENV_STREAM_FPS)s",STREAM_BITRATE="%(ENV_STREAM_BITRATE)s",RTMP_URL="%(ENV_RTMP_URL)s",YOUTUBE_STREAM_KEY="%(ENV_YOUTUBE_STREAM_KEY)s",SNAPSHOT_DIR="%(ENV_SNAPSHOT_DIR)s",SNAPSHOT_INTERVAL="%(ENV_SNAPSHOT_INTERVAL)s"
stdout_logfile=/var/log/supervisor/stream.log
stderr_logfile=/var/log/supervisor/stream.err

[program:bot]
command=python3 /app/bot.py
autorestart=true
priority=100
environment=DISPLAY=":1",YOUTUBE_API_KEY="%(ENV_YOUTUBE_API_KEY)s",YOUTUBE_VIDEO_ID="%(ENV_YOUTUBE_VIDEO_ID)s",ROBLOX_GAME_ID="%(ENV_ROBLOX_GAME_ID)s"
stdout_logfile=/var/log/supervisor/bot.log
stderr_logfile=/var/log/supervisor/bot.err
